# autostart/conky.desktop
[Desktop Entry]
Type=Application
Exec=sh -c "~/.config/conky/start.sh"
Name=Conky
Comment=Autostart conky at login
